# data/greeter-autostart/polkit-mate-authentication-agent-1.desktop.in
[Desktop Entry]
Type=Application
_Name=PolicyKit Authentication Agent
Exec=/usr/lib/policykit-1-mate/polkit-mate-authentication-agent-1
OnlyShowIn=MATE;
